# .config/autostart/pika-monitor.desktop
[Desktop Entry]
Name=Pika Backup Monitor
Type=Application
Exec=flatpak run --command=pika-backup-monitor org.gnome.World.PikaBackup
NoDisplay=true
